# hn-gopher/etc/cron.d/hngopher
0 * * * * root /usr/local/bin/scrape_hn --log-file /var/log/hngopher
