# hn-gopher/etc/cron.d/hngopher
0 * * * *   /usr/local/bin/scrape_hn --log-file /var/log/scrape_hn