# hn-gopher/etc/cron.d/hngopher
0 * * * * root /usr/local/bin/hn-scrape --log-file /var/log/hngopher
30 23 * * * root /usr/local/bin/hn-archive --log-file /var/log/hngopher
